# systemd/service/mc-backup.service
[Unit]
Description=Minecraft Server - Backup (service)
Wants=network-online.target
Requires=mc-backup.timer
Before=mc-backup-upload.service

[Service]
User=minecraft
Group=minecraft
Type=simple
ExecStart=scripts/backup.sh
# SET REAL PATH
Nice=8

[Install]
WantedBy=multi-user.target
